# config/frpc.supervisor.conf
[program:frpc_vnc]
command=/root/center_frpc_vnc/bin/frpc 49997
stopasgroup=true
killasgroup=true
restart=always
user=root
