# config/frpc.supervisor.ws.banner.conf
[program:frpc_ws_banner]
command=/root/center_frpc_ws/bin/frpc banner 49999
stopasgroup=true
killasgroup=true
restart=always
user=root
